# backups/.crontab

0 3 * * 1,5 /usr/bin/python /home/ericp/backups/backup.py